# deploy/vetinari.service
# =============================================================================
# EXAMPLE SCAFFOLDING - review before use in production.
#
# - User/Group are set to "vetinari"; create this account or substitute your
#   service account before enabling.
# - amw-kernel-server must be installed on PATH or replace ExecStart with the
#   absolute path used by your package.
# - ReadWritePaths lists the directories Vetinari writes to at runtime.
#   ProtectSystem=strict blocks all other writes.
# =============================================================================

[Unit]
Description=Vetinari AM Workbench Kernel
After=network-online.target
Wants=network-online.target

[Service]
# Type=simple: the kernel does not emit sd_notify readiness; Type=notify would
# cause systemd to wait forever for a notification that never arrives.
# Release readiness gates must query http://127.0.0.1:${AMW_KERNEL_PORT}/ready.
Type=simple
User=vetinari
Group=vetinari
WorkingDirectory=/opt/vetinari
ExecStart=/usr/local/bin/amw-kernel-server --host ${AMW_KERNEL_HOST} --port ${AMW_KERNEL_PORT}
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
# Allow writes to the directories Vetinari uses at runtime
ReadWritePaths=/opt/vetinari/data

# Resource limits
MemoryMax=4G
CPUQuota=200%

# Canonical Rust kernel environment variables
Environment=AMW_KERNEL_HOST=127.0.0.1
Environment=AMW_KERNEL_PORT=8080
Environment=AMW_KERNEL_STATE_DIR=/opt/vetinari/data

[Install]
WantedBy=multi-user.target
